# roles/nftables/templates/refresh-discord-ips.service.j2
[Unit]
Description=Refresh Discord IPs in nftables
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/refresh-discord-ips.sh
